# templates/default/atlantis.upstart.conf.erb
description "Atlantis"

start on runlevel [2345]
stop on runlevel [06]

respawn
respawn limit 10 10

umask 022

chdir <%= @vars[:atlantis_home] || '/opt/atlantis' %>
env GOMAXPROCS="2"
env PATH="/usr/local/bin:/usr/bin:/bin:/sbin"

kill signal TERM
kill timeout <%= @vars[:timeout_stop_sec] || 5 %>
reload signal HUP

script
  # debug
  set -x

  # start the service
  exec start-stop-daemon --start --chuid <%= @vars[:atlantis_user] %> --exec <%= @vars[:atlantis_bin_location] %> -- server --config <%= @vars[:atlantis_home] %>/<%= @vars[:atlantis_config_name] %>
end script
